# packaging/systemd/cpu-alert.service
[Unit]
Description=CPU Alert Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cpu-alert -config /etc/cpu-alert/config.yaml
WorkingDirectory=/etc/cpu-alert
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
